# data/export/systemd/process.service.erb
[Unit]
PartOf=<%= service.full_name %>.target

[Service]
WorkingDirectory=<%= work_directory %>
Environment=XDG_RUNTIME_DIR=/run/user/%I
Environment=PORT=%i
Environment=KARMA_IDENTIFIER=%n
ExecStart=/bin/bash -lc '<%= service.command %>'
Restart=always
StandardInput=null
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n
KillMode=mixed
TimeoutStopSec=<%= service.timeout_stop %>
<% if service.class.config_cpu_quota.present? %>CPUAccounting=true<% end %>
<% if service.class.config_cpu_quota.present? %>CPUQuota=<%= service.class.config_cpu_quota %>%<% end %>
<% if service.class.config_memory_max.present? %>MemoryAccounting=true<% end %>
<% if service.class.config_memory_max.present? %>MemoryLimit=<%= service.class.config_memory_max %><% end %>
