# data/export/systemd/process.service.erb
[Unit]
PartOf=<%= service.full_name %>.target

[Service]
WorkingDirectory=<%= work_directory %>
Environment=PORT=%i
Environment=KARMA_IDENTIFIER=%n
Environment=KARMA_ENV=<%= Karma.env %>
Environment=KARMA_PROJECT_ID=<%= Karma.karma_project_id %>
Environment=KARMA_USER_ID=<%= Karma.karma_user_id %>
Environment=KARMA_AWS_USER_ACCESS_KEY=<%= Karma.aws_access_key_id %>
Environment=KARMA_AWS_USER_SECRET_ACCESS_KEY=<%= Karma.aws_secret_access_key %>
ExecStart=/bin/bash -lc '<%= service.config_command %>'
Restart=<%= service.config_auto_restart ? 'always' : 'no' %>
StandardInput=null
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n
KillMode=mixed
TimeoutStopSec=<%= service.config_timeout_stop %>
<% if service.config_cpu_quota.present? %>CPUAccounting=true<% end %>
<% if service.config_cpu_quota.present? %>CPUQuota=<%= service.config_cpu_quota %>%<% end %>
<% if service.config_memory_max.present? %>MemoryAccounting=true<% end %>
<% if service.config_memory_max.present? %>MemoryLimit=<%= service.config_memory_max %><% end %>
