# data/export/systemd/process.service.erb
[Unit]
PartOf=<%= project_name %>-<%= service.name %>.target

[Service]
WorkingDirectory=<%= work_directory %>
Environment=PORT=%i
ExecStart=/bin/bash -lc '<%= service.command %>'
Restart=always
StandardInput=null
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n
KillMode=mixed
TimeoutStopSec=5
<% if service.class.config_cpu_quota.present? %>CPUAccounting=true<% end %>
<% if service.class.config_cpu_quota.present? %>CPUQuota=<%= service.class.config_cpu_quota %>%<% end %>
<% if service.class.config_memory_max.present? %>MemoryMax=<%= service.class.config_memory_max %><% end %>
