# data/export/systemd/process_master.target.erb
[Unit]
PartOf=<%= project_name %>.target
